# scripts/com.agentra.hydra.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentra.hydra</string>

    <key>Program</key>
    <string>/Users/omoshola/.local/bin/hydra-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/omoshola/.local/bin/hydra-server</string>
        <string>--port</string>
        <string>7777</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/omoshola/Library/Logs/hydra-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/omoshola/Library/Logs/hydra-server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/omoshola</string>
    </dict>
</dict>
</plist>
